# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/storage/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:web]
command=php artisan serve --host=0.0.0.0 --port=%(ENV_PORT)s
directory=/var/www
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/web.log
stderr_logfile=/var/www/storage/logs/web.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:queue]
command=php artisan queue:work database --daemon --timeout=300 --tries=1 --sleep=3 --max-jobs=100 --max-time=3600
directory=/var/www
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/queue.log
stderr_logfile=/var/www/storage/logs/queue.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
